# src/main/resources/kafka-three-supervisord.ini
[program:kafka-three]
command=java -server -XX:+UseG1GC -DCONFIG_FILE=/data/kafka/test.properties -jar /data/kafka/kafka-three.jar
process_name=%(program_name)s
numprocs=1
directory=/data/kafka/
autostart=false
autorestart=false
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=30
user=cagricelebi
redirect_stderr=true
stdout_logfile=/data/kafka/out.log